# tree/etc/cron.daily/00logwatch
#!/bin/bash

# Do nothing.
# This is provided by sysconf.nef.monitor to disable default daily logwatch
# and enable own scripts instead.
